# conf/cron.d/wptpl-cron-job.example
*/5 * * * *  www-data	/usr/bin/php -f WP_ROOT/wp-cron.php DOING_CRON > /dev/null 2>&1
